# openvpn-as/root/etc/s6-overlay/s6-rc.d/svc-openvpn/finish
#!/usr/bin/with-contenv bash

kill "$(cat /openvpn/pid/openvpn.pid)"
/usr/local/openvpn_as/scripts/openvpn_service_cleanup
