# openvpn-as/root/etc/s6-overlay/s6-rc.d/svc-openvpn/run
#!/usr/bin/with-contenv bash

/usr/local/openvpn_as/scripts/openvpnas --nodaemon --pidfile=/openvpn/pid/openvpn.pid --logfile=/config/log/openvpn.log
